# roles/flask_supervisor_example/templates/flask.conf.j2
[program:{{flask_app}}]
# Run as our user
user={{ ansible_user_id }}

# Run the templated script
command={{ flask_directory }}/{{flask_app}}.sh
